# deploy/systemd/opencode-router.service.example
[Unit]
Description=OpenCode Tailnet Router
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/opencode-router
ExecStart=/usr/bin/node /opt/opencode-router/router/vps-opencode-router.js
Environment=OPENCODE_ROUTER_CACHE_DIR=/opt/opencode-router/cache
Restart=always
RestartSec=2
Environment=OPENCODE_ROUTER_HOST=127.0.0.1
Environment=OPENCODE_ROUTER_PORT=33102
Environment=OPENCODE_ROUTER_INSPECT_TIMEOUT_MS=8000
Environment=OPENCODE_ROUTER_WARM_TIMEOUT_MS=30000
Environment=OPENCODE_ROUTER_HTML_TIMEOUT_MS=8000
Environment=OPENCODE_ROUTER_RECOVERY_HTML_TIMEOUT_MS=15000

[Install]
WantedBy=multi-user.target
